# dev-java/sonar-bin/files/sonar-init-32
#!/sbin/runscript

depend() {
    need net
    use dns logger
}

RUN_AS_USER=sonar

start() {
	ebegin "Starting ${SVCNAME}"
	su $RUN_AS_USER -c "/opt/sonar/bin/linux-x86-32/sonar.sh start"
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	su $RUN_AS_USER -c "/opt/sonar/bin/linux-x86-32/sonar.sh stop"
	eend $?
}
